# crontab
* * * * * root ./back-up.sh --type full
